# deploy/systemd/saga.service
[Unit]
Description=Saga service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/saga
EnvironmentFile=/etc/saga/saga.env
ExecStart=/usr/bin/env bun run --cwd /opt/saga --filter @saga/service start
Restart=on-failure
RestartSec=5
User=saga
Group=saga

[Install]
WantedBy=multi-user.target
