# deploy/systemd/saga.service
[Unit]
Description=Saga service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/saga/apps/service
EnvironmentFile=/etc/saga/saga.env
# Direct exec so the service process receives systemd's SIGTERM itself;
# `bun run --filter` supervision does not forward SIGTERM to the script child.
ExecStart=/usr/bin/env node --import tsx src/main.ts
Restart=on-failure
RestartSec=5
User=saga
Group=saga

[Install]
WantedBy=multi-user.target
